# rootfs/etc/s6-overlay/s6-rc.d/tcp-ws/run
#!/bin/sh
#exec 2>&1
#exec python3 /app/server.py 2>&1
exec python3 /app/server.py 2>&1
