# rootfs/etc/s6-overlay/s6-rc.d/tcp-ws-log-prepare/up
if { mkdir -p /var/log/tcp-ws }
if { chown nobody:nogroup /var/log/tcp-ws }
chmod 02755 /var/log/tcp-ws